# s6/alpine-root/etc/s6-overlay/s6-rc.d/_startup/up
#!/command/execlineb
foreground { bash -e /usr/local/bin/_pre_startup.sh }